# configs/autostart/spotify.desktop
[Desktop Entry]
Type=Application
Name=spotify
Exec=spotify
